# confs/gatekeeper_supervisor.conf
[program:Gatekeeper_Server]
command=/usr/bin/python /root/gatekeeper_webserv.py
stderr_logfile=/var/log/supervisor/err_supervisor.log
user=root